# systemd/.config/systemd/user/dunst.service
[Unit]
Description=Dunst Notification Daemon

[Service]
ExecStart=/usr/bin/dunst
Restart=on-failure

[Install]
WantedBy=default.target

