# systemd/.config/systemd/user/wallpaper-changer.service

[Unit]
Description=Change wallpaper every 3 hours

[Service]
Type=oneshot
ExecStart=/home/pratik/.config/scripts/wallpaper-changer.sh
